# conf/nrest-api.service
[Unit]
Description=Workflow Templates API Server
After=network.target
After=sqlite3.service

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/workflow-templates
ExecStart=/opt/workflow-templates/nrest-api
Environment="PORT=5680"
Environment="DATABASE_FILE=/opt/workflow-templates/workflow_templates.db"
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
KillSignal=SIGTERM

# Security hardening
PrivateTmp=true
ProtectSystem=full
NoNewPrivileges=true
PrivateDevices=true

[Install]
WantedBy=multi-user.target
